# deploy/andrewmccall-backend.service
[Unit]
Description=AndrewMcCall backend (Rust/hyper)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
# WorkingDirectory is the backend source dir so the binary finds ./migrations,
# writes ./logs, and loads ./.env (dotenv). EnvironmentFile injects the same
# vars explicitly, so it works even if the working directory ever changes.
# Adjust these paths if the repo lives somewhere other than /root/AndrewMcCall.
WorkingDirectory=/root/AndrewMcCall/backend
EnvironmentFile=/root/AndrewMcCall/backend/.env
# Bind to loopback only — nginx proxies to it. Drop ADDR to bind all interfaces.
Environment=ADDR=127.0.0.1
Environment=PORT=3000
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/andrewmccall-backend
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
